# supervisord.conf
[supervisord]
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor
nodaemon = true

[program:nginx]
command = /usr/sbin/nginx
stdout_logfile = /var/log/supervisor/%(program_name).log
stderr_logfile = /var/log/supervisor/%(program_name).log
autorestart = true

[program:influxdb]
command = /usr/bin/influxdb -config=config.toml
stdout_logfile = /var/log/supervisor/%(program_name).log
stderr_logfile = /var/log/supervisor/%(program_name).log
autorestart = true